# python/supervisor/doorbot-check.conf

[program:doorbot-check]
command=/opt/doorbot/check.sh
