# python/supervisor/doorbot.conf

[program:doorbot]
command=/opt/doorbot/run.sh
